# dynamic_lock.service
# dynamic_lock.service — systemd user service for bluetooth proximity lock
#
# installed by install.sh to ~/.config/systemd/user/
# manage with:
#   systemctl --user status dynamic_lock
#   systemctl --user restart dynamic_lock
#   systemctl --user stop dynamic_lock

[Unit]
Description=Bluetooth Dynamic Lock
After=bluetooth.target graphical-session.target
Wants=bluetooth.target graphical-session.target
StartLimitIntervalSec=120
StartLimitBurst=4

[Service]
Type=simple
ExecStart=%h/.local/bin/dynamic_lock.sh

# fast shutdown — script traps SIGTERM and exits immediately,
# but this acts as a safety net in case it doesn't
TimeoutStopSec=5
KillSignal=SIGTERM

# dbus session address — needed for loginctl lock-session and notify-send
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

# restart on crash, not on clean exit
Restart=on-failure
RestartSec=20

# logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
